# bbcguild-raffle.service
[Unit]
Description=BBC Guild Raffle Web Application
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=bbcguild-raffle
WorkingDirectory=/home/technitaur/win.bbcguild.com
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/local/bin/gunicorn --config gunicorn.conf.py wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5s
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
